# package/network/services/stunnel/files/stunnel.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2008 OpenWrt.org

START=90
RUN_D=/var
PID_F=$RUN_D/stunnel.pid

start() {
	[ ! -f $PID_F ] && stunnel
}

stop() {
	[ -f $PID_F ] && {
		kill $(cat $PID_F)
		kill -9 $(cat $PID_F)
		rm -f $PID_F
	}
}
