# package/base-files/files/etc/uci-defaults/13_enable_services
#!/bin/sh

/etc/init.d/aria2 enable
/etc/init.d/samba enable
/etc/init.d/fstab enable

exit 0

